# scripts/templates/pumaguard.service.j2
[Unit]
Description=PumaGuard Wildlife Monitoring Service
Documentation=http://pumaguard.rtfd.io/
After=network-online.target dnsmasq.service
Wants=network-online.target

[Service]
Type=simple
User={{ pumaguard_user | default(ansible_user) }}
Group={{ pumaguard_group | default(ansible_user) }}
WorkingDirectory={{ pumaguard_watch_dir | default('/home/' + ansible_user + '/pumaguard') }}

# Set XDG Base Directory environment variables for proper config/cache/data locations
Environment="HOME=/home/{{ pumaguard_user | default(ansible_user) }}"
Environment="XDG_CONFIG_HOME=/home/{{ pumaguard_user | default(ansible_user) }}/.config"
Environment="XDG_CACHE_HOME=/home/{{ pumaguard_user | default(ansible_user) }}/.cache"
Environment="XDG_DATA_HOME=/home/{{ pumaguard_user | default(ansible_user) }}/.local/share"

# Set PulseAudio/PipeWire environment variables for audio playback
Environment="XDG_RUNTIME_DIR=/run/user/{{ pumaguard_uid | default('1000') }}"
Environment="PULSE_RUNTIME_PATH=/run/user/{{ pumaguard_uid | default('1000') }}/pulse"
Environment="DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/{{ pumaguard_uid | default('1000') }}/bus"

ExecStart=/opt/pumaguard/bin/pumaguard server --debug {{ pumaguard_watch_dir | default('/home/' + ansible_user + '/pumaguard') }}
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
